# templates/rkt-image-gc.timer.j2
[Unit]
Description=Run rkt Image Garbage Collector {{rkt_image_gc_frequency}}

[Timer]
OnCalendar={{rkt_image_gc_frequency}}
Persistent=true

[Install]
WantedBy=timers.target
